# app1-API/crontab
* * * * * python3.8 /app/app.py >> /tmp/pylog.txt

